# node/manifest/config/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/local/bin/init_node.sh
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0
